# crontab
0 14 * * * cd /app && npm start >> /var/log/cron.log 2>&1
